# packages/pyrefdev/pyrefdev-2025.11.tar.gz/pyrefdev-2025.11/pyrefdev.service
[Unit]
Description=Web server for pyref.dev
After=network.target

[Service]
User=root
WorkingDirectory=/root/pyref.dev
ExecStart=uv run uvicorn pyrefdev.server:app --host=0.0.0.0 --port=80
Restart=always
RestartSec=5
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=pyref.dev

[Install]
WantedBy=multi-user.target
